# supervisor-simplehttpd.ini
[program:httpd]
command=/usr/bin/python -m SimpleHTTPServer 8080
directory=/home/git/repos
user=git
redirect_stderr=true
